# fstore/initrc/rc.d/fstored.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: fstored
# REQUIRE: LOGIN

. /etc/rc.subr

name="fstored"
rcvar="fstored_enable"

pidfile="@srv_rundir@"/fstored.pid
command="@prefix@/sbin/${name}"
command_args="-daemon"
procname="@prefix@/sbin/${name}"

load_rc_config ${name}

: ${fstored_enable:="NO"}

run_rc_command "$1"
#EOF
